# deploy/dolmyanfit.service
[Unit]
Description=DolmyanFit Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dolmyanfit
ExecStart=/root/dolmyanfit/.venv/bin/python /root/dolmyanfit/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
